# systemd/system/ti-uim.service
[Unit]
Description=User Mode Init Manager for TI shared transport

[Service]
ExecStartPre=-/sbin/modprobe -q btwilink
ExecStart=/usr/sbin/uim
